# debian/uchroma.postinst
#!/bin/sh

set -e

if [ "$1" = configure ]; then

    # update udev with new configuration
    udevadm control --reload-rules
    udevadm trigger -s usb -s hidraw -c add -a idVendor=1532

    # enable uchromad for all users
    if which systemctl >/dev/null 2>&1; then
        systemctl --global enable uchromad.service || true
    fi
fi

exit 0
